# private_dot_config/systemd/user/tmux.service.tmpl
[Unit]
Description=Start tmux session named "zero" on login

[Service]
ExecStart=%h/.config/tmux/start_tmux.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=default.target
